# ci/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Jimm CI
After=network-online.target jimm-ci.service
Wants=network-online.target

[Service]
Type=simple
User=gh-runner
Group=gh-runner
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/gh-runner/.cloudflared/config.yml run jimm-ci
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/gh-runner/.cloudflared
PrivateTmp=true

[Install]
WantedBy=multi-user.target
