# ci/systemd/jimm-ci-nightly.service
[Unit]
Description=Jimm.jl nightly full parity sweep
After=jimm-ci.service

[Service]
Type=oneshot
User=gh-runner
Group=gh-runner
WorkingDirectory=/opt/jimm-ci/Jimm.jl/ci

Environment=JIMM_CI_STATE_DIR=/var/lib/jimm-ci
Environment=JIMM_CI_APP_ID=REPLACE_WITH_APP_ID
Environment=JIMM_CI_INSTALLATION_ID=REPLACE_WITH_INSTALLATION_ID
Environment=JIMM_CI_PRIVATE_KEY_FILE=/etc/jimm-ci/app.pem
Environment=JIMM_CI_WEBHOOK_SECRET_FILE=/etc/jimm-ci/webhook-secret
Environment=JIMM_CI_HF_TOKEN_FILE=/etc/jimm-ci/hf-token
Environment=JIMM_CI_REPO_OWNER=REPLACE_WITH_OWNER
Environment=JIMM_CI_REPO_NAME=Jimm.jl
Environment=JIMM_CI_PUBLIC_URL=https://jimm-ci.REPLACE_WITH_DOMAIN
Environment=JIMM_CI_JULIA=/usr/local/bin/julia
Environment=JULIA_NUM_THREADS=4

ExecStart=/usr/bin/env uv run --project /opt/jimm-ci/Jimm.jl/ci jimm-ci-nightly

# This unit does its own port binding (queue + worker), so it cannot run while
# the listener service holds the same git mirror. Stop the listener first if
# you want absolute isolation; otherwise the bare mirror and depot are shared
# safely (git fetch is atomic, Pkg.test is read-only on the depot).
